# dbus-demo/code/dbus-1.12.16/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/home/ubuntu/dbus-demo/lib/dbus/var/run/dbus/system_bus_socket
